# bedtime.timer
[Unit]
Description=Run bedtime-mode at 8 PM (2 hours before bedtime)

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
