# wakeup-monitor.service
[Unit]
Description=Wake up the monitor and reset RGB

[Service]
Type=oneshot
ExecStart=/home/benjamin/256/morning-wakeup.sh
